# supervisor.conf
[program:modbus]
user = root
command = /usr/bin/python /app/app.py
autostart = true
autorestart = true
startsecs = 5
